# Caldera/extraFiles/target_root/etc/profile.d/fluke_cda_profile.sh

export PRODUCT_NAME=Caldera
PATH=$PATH:/home/$PRODUCT_NAME/bin

export LS_OPTS='--color=auto'
alias ls='ls ${LS_OPTS}'

export PRODUCT_DIR=/home/$PRODUCT_NAME/

