# Nighthawk/extraFiles/etc/profile.d/fluke_cda_profile.sh

PATH=$PATH:/home/Nighthawk/bin

export LS_OPTS='--color=auto'
alias ls='ls ${LS_OPTS}'

export PRODUCT_NAME=Nighthawk
export PRODUCT_DIR=/home/Nighthawk

