# Caldera/target/var/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/usr/sbin/dropbearmulti
/etc/dropbear/
/usr/sbin/dropbearkey
/lib/systemd/system/dropbearkey.service
/usr/sbin/dropbearconvert
/lib/systemd/system/dropbear.socket
/usr/sbin/dropbear
/etc/pam.d/dropbear
/lib/systemd/system/dropbear@.service
